# osint/Update.cron
15     3     *     *     * ./GrabAndPush.sh
